# muraai-contracts.service
[Unit]
Description=Muraai Contract Renewal Reminder App
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/muraai-contract-reminder
Environment=PATH=/path/to/muraai-contract-reminder/venv/bin
EnvironmentFile=/path/to/muraai-contract-reminder/.env
ExecStart=/path/to/muraai-contract-reminder/venv/bin/gunicorn -w 4 -b 127.0.0.1:5000 src.main:app
Restart=on-failure
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
